# crontab
15 * * * *  rsync -a /mnt/bucket /mnt/logfiles_incoming
 0 0 0 * *  /mnt/bucket/appengine.googleapis.com/request_log/`date +"%Y"`/`date +"%m"`/`date +"%d"`
