# crontab
15 * * * *  rsync -a /mnt/bucket /mnt/logfiles_incoming